# Disco/lte/lib/70-huawei-e3372h-153.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="12d1", ATTRS{idProduct}=="157d", RUN+="/data/ftp/internal_000/lte/bin/lte-hotplug"
